# debian/prerm
#!/bin/sh

HAPP=houseportal

case $1 in
    remove|purge)
        if [ -e /lib/systemd/system/${HAPP}.service ] ; then systemctl stop ${HAPP} ; systemctl disable ${HAPP} ; fi
        ;;
esac

